# crontab
SHELL=/bin/sh
0-59 * * * * manage.py check



